# private_dot_config/systemd/user/hyprland-keep-ncspot-alive.service
[Unit]
Description=Hyprland "keep ncspot alive"

[Service]
Type=oneshot
ExecStart=hyprctl dispatch exec %h/.local/bin/hyprland-keep-ncspot-alive.sh

