# private_dot_config/systemd/user/hyprland-laptop-controller.service
[Unit]
Description=Hyprland Laptop Controller

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprland-laptop-controller.sh

