# private_dot_config/systemd/user/laptop-controller.service
[Unit]
Description=Hyprland Laptop Display Controller

[Service]
Type=oneshot
ExecStart=%h/.local/bin/laptop-controller.sh

