# private_dot_config/systemd/user/lid-robot.service
[Unit]
Description=lid-robot

[Service]
Type=oneshot
ExecStart=/home/chris/scripts/lid-robot.sh
